# systemd/market-data-analyze.service
[Unit]
Description=Market Data Analysis Scheduler (9:00 AM ET)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/md
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m src.main schedule-analyze
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/projects/md/logs/analyze-scheduler.log
StandardError=append:/home/ubuntu/projects/md/logs/analyze-scheduler.error.log

[Install]
WantedBy=multi-user.target
