# systemd/market-data-fetch.service
[Unit]
Description=Market Data Fetch Scheduler (4:15 PM ET)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/md
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m src.main schedule-fetch
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/projects/md/logs/fetch-scheduler.log
StandardError=append:/home/ubuntu/projects/md/logs/fetch-scheduler.error.log

[Install]
WantedBy=multi-user.target
